# services/sports-watchdog.service
[Unit]
Description=Sports Display Watchdog Monitor
After=network.target sports-display.service

[Service]
Type=simple
User={{USER}}
WorkingDirectory={{PROJECT_DIR}}
Environment="PATH={{VENV_PATH}}/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart={{VENV_PATH}}/bin/python src/watchdog_monitor.py

# Restart policy
Restart=always
RestartSec=60

# This is a simple monitoring service, no watchdog needed
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target